# deploy/hf-newsletter.service
[Unit]
Description=Hugging Face Newsletter Service
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/hf_newsletter
Environment=PYTHONPATH=/root/hf_newsletter
ExecStart=/root/hf_newsletter/venv/bin/python3 src/cli.py
StandardOutput=append:/root/hf_newsletter/logs/service.log
StandardError=append:/root/hf_newsletter/logs/service.log

[Install]
WantedBy=multi-user.target
